# content/s6-overlay/s6-rc.d/init-sieve/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -e

RSPAMD_HOST=${RSPAMD_HOST:-"rspamd-svc"}

echo "Refreshing /data/sieve"
rm -rf /data/sieve/*
cp -r /configfiles/sieve/* /data/sieve

echo "Compiling sieve files"
for f in /data/sieve/*.sieve; do
  [ -f "$f" ] || continue  # Skip if file doesn't exist
  echo "Compiling sieve file $f"
  sievec "$f" || echo "Warning: Failed to compile $f"
done

echo "Setting up rspamd_learn_ham.sh and rspamd_learn_spam.sh to use RSPAMD_HOST: $RSPAMD_HOST"
for script in /data/sieve/rspamd_learn_ham.sh /data/sieve/rspamd_learn_spam.sh; do
  [ -f "$script" ] && sed -i "s|RSPAMD_HOST|${RSPAMD_HOST}|g" "$script"
done

echo "Changing user:group to vmail:mail for /data/sieve"
chown -R vmail:mail /data/sieve

echo "Making sieve binaries and scripts executable"
find /data/sieve -name "*.svbin" -exec chmod +x {} \;
find /data/sieve -name "*.sh" -exec chmod +x {} \;

echo "Sieve initialization complete"
